# dist/llamaloops-supervisor.conf
[program:llamaloops]
directory=/home/dulan/flask-app
command=/home/dulan/flask-app/venv/bin/gunicorn -w 3 run:app
user=dulan
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/flask-app/llamaloops.err.log
stdout_logfile=/var/log/flask-app/llamaloops.out.log
