# dist/llamaloops-supervisor.conf
[program:audimax]
directory=/home/iamarobot64/audimax/dist
command=/home/iamarobot64/audimax/dist/venv/bin/gunicorn -w 3 run:app
user=iamarobot64
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/supervisor/audimax.err.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/audimax.out.log
stdout_logfile_maxbytes=10MB